# ubuntu-gui-base/setids.service
[Unit]
Description=Set_IDs

[Service]
ExecStart=/usr/local/bin/setids.sh

[Install]
WantedBy=multi-user.target